# deploy/systemd/rc-fetch-secrets.service
[Unit]
Description=Fetch Rising Constellation env vars from AWS Secrets Manager
Documentation=https://docs.aws.amazon.com/secretsmanager/latest/userguide/
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Runs as root because it writes to /etc/rc/env (mode 0600, owned by rc:rc).
# The actual API call uses the EC2 instance role's credentials via IMDS;
# no AWS keys are needed on the host.
ExecStart=/usr/local/bin/rc-fetch-secrets

# If Secrets Manager is unreachable, fail fast — rc.service will not start
# without the env file.
TimeoutStartSec=60s

[Install]
WantedBy=rc.service
